# opc-sim/opc-simulator.service
[Unit]
Description=OPC UA Simulator — mintpower lab
After=network.target

[Service]
Type=simple
User=mintpower
WorkingDirectory=/home/mintpower/lab/k3s/opc-sim
ExecStart=/home/mintpower/lab/k3s/opc-sim/venv/bin/python3 opc_server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
